# deploy/parakeet-daemon.service
[Unit]
Description=Parakeet STT Daemon (uv run)
After=network.target sound.target

[Service]
# UPDATE THIS PATH
WorkingDirectory=/path/to/parakeet-stt/parakeet-stt-daemon
ExecStart=/usr/bin/env uv run parakeet-stt-daemon
Environment=PARAKEET_DEVICE=cuda
Environment=PARAKEET_MIC_DEVICE=pulse
Environment=PARAKEET_STREAMING_ENABLED=true
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
