# deploy/parakeet-ptt.service
[Unit]
Description=Parakeet PTT Client
After=graphical-session.target network.target
Wants=graphical-session.target

[Service]
WorkingDirectory=/home/hugo/Documents/Engineering/parakeet-stt/parakeet-ptt
ExecStart=/home/hugo/Documents/Engineering/parakeet-stt/parakeet-ptt/target/release/parakeet-ptt
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
